# crontab
* * * * * echo $APP_URL >> /var/log/cron.log 2>&1
#end line